# service_files/chassis-poweroff.target
[Unit]
Description=Chassis (Power Off)
After=multi-user.target
Conflicts=chassis-poweron.target